# .config/systemd/user/i3battcheck.timer
[Unit]
Description=Run i3battcheck every minute

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
