# .config/systemd/user/redshift.service
[Unit]
Description=Redshift

[Service]
Type=simple
ExecStart=/usr/bin/redshift-gtk -l geoclue2

[Install]
WantedBy=default.target
